# bezahlonlinekassa/DEBIAN/postinst
cp /home/kassa/getty@tty1.service /etc/systemd/system/getty.target.wants
mv /usr/share/nginx/html/index.html /usr/share/nginx/html/index.html.orig
cp /usr/share/bezahlonline/register/html/assets /usr/share/nginx/html/
cp /usr/share/bezahlonline/register/html/index.html /usr/share/nginx/html/
systemctl restart nginx
systemctl daemon-reload
useradd -r -G dialout kassa >/dev/null 2>&1

ZTV_LOG_DIR=/var/log/zvt
if [ ! -d ${ZTV_LOG_DIR} ]; then
    mkdir ${ZTV_LOG_DIR}
fi
chown kassa /var/log/zvt

systemctl enable gm65
systemctl enable ptapi
systemctl enable printapi
systemctl enable rfid
systemctl start gm65
systemctl start ptapi
systemctl start printapi
systemctl start rfid

