# bezahlonlinekassa/DEBIAN/postinst
cp /home/kassa/getty@tty1.service /etc/systemd/system/getty.target.wants
mv /usr/share/nginx/html/index.html /usr/share/nginx/html/index.html.orig
cp -r /usr/share/bezahlonline/register/html/assets /usr/share/nginx/html/
cp /usr/share/bezahlonline/register/html/index.html /usr/share/nginx/html/

systemctl restart nginx
systemctl daemon-reload

# create user 'kassa'
USER=kassa
useradd -s /bin/bash -r -G dialout $USER >/dev/null 2>&1

ZTV_LOG_DIR=/var/log/zvt
if [ ! -d ${ZTV_LOG_DIR} ]; then
    mkdir ${ZTV_LOG_DIR}
fi
chown $USER:$USER /var/log/zvt
chown -R $USER:$USER /home/$USER
cp /home/$USER/getty@tty1.service /etc/systemd/system/getty.target.wants/

# won't work
# certutil -d sql:/home/$USER/.pki/nssdb -A -t "C,," -n bezahl-online -i /home/$USER/myCA.pem

systemctl enable gm65
systemctl enable ptapi
systemctl enable printapi
systemctl enable rfid
systemctl start gm65
systemctl start ptapi
systemctl start printapi
systemctl start rfid

# this should login the kassa and start it up
systemctl restart getty@tty1 
