# bezahlonlinekassa/DEBIAN/prerm
systemctl stop nginx
systemctl stop gm65
systemctl stop ptapi
systemctl stop printapi
systemctl stop rfid
systemctl disable gm65
systemctl disable ptapi
systemctl disable printapi
systemctl disable rfid

# stop X server - won't remove it though
killall Xorg